# lib/rae-dvdhrm/files/etc-systemd-user.sway.service
[Unit]
Before=sway-session.target
Description=Sway Compositor
PartOf=sway-session.target
RefuseManualStart=yes
RefuseManualStop=yes

[Service]
Environment=XDG_CURRENT_DESKTOP=sway
ExecStart=/usr/bin/sway
ExecStopPost=/usr/bin/systemctl --user unset-environment DISPLAY SWAYSOCK WAYLAND_DISPLAY XDG_CURRENT_DESKTOP
NotifyAccess=all
Type=notify
