# lib/rae-dvdhrm/files/etc-systemd-user.sway-session.target
[Unit]
Description=Sway Session Target

BindsTo=graphical-session.target
Before=graphical-session.target

After=sway.service
BindsTo=sway.service
